# heating-server.service
[Unit]
Description=Smart Heating Server

[Service]
ExecStart=/home/pi/raspi-heating-control/server.sh

[Install]
WantedBy=multi-user.target
