# heating-checker.service
[Unit]
Description=Smart Heating Checker

[Service]
ExecStart=/home/pi/raspi-heating-control/checker.py

[Install]
WantedBy=multi-user.target
